# deploy/docker/s6-rc.d/valkey/run
#!/command/execlineb -P
# Valkey service runner for s6-overlay

foreground { echo "Starting Valkey..." }

# Run as valkey user (or root if user doesn't exist)
valkey-server /etc/valkey/valkey.conf
