# deploy/docker/s6-rc.d/vault-init/up
/bin/bash -c "
export VAULT_ADDR=http://127.0.0.1:8200
echo 'Waiting for Vault to be ready...'
attempts=0
while ! curl -s http://127.0.0.1:8200/v1/sys/health > /dev/null 2>&1; do
    attempts=$((attempts + 1))
    if [ $attempts -gt 30 ]; then
        echo 'ERROR: Vault not responding after 30 seconds'
        exit 1
    fi
    sleep 1
done
echo 'Vault is listening, running unseal script...'
exec /opt/mira/app/deploy/docker/scripts/vault-unseal.sh
"
